# backup/offsite.service
[Unit]
Description=Offsite Backup

[Service]
Type=oneshot
ExecStart=%h/.local/bin/offsite
Restart=no
